# systemd/user/macro-keyboards.target
[Unit]
Description=All Macro Keyboards
Documentation=http://berry.local/code/crusty-desktop/crusty-keyboard

[Install]
WantedBy=default.target
